# etc/cron.hourly/cron-drupal-sync-kmd-data
#!/bin/bash

cd /var/www/w3.uib.no && bin/site-drush uib-area-import-kmd-news
