# etc/cron.daily/cron-drupal-sync-w3-logs
#!/bin/bash
cd /var/www/w3.uib.no
su w3-drupal -c 'bin/sync-w3-logs'
