# etc/cron.daily/cron-drupal-sync-w3-logs
#!/bin/bash
cd /var/www/w3.uib.no
log="/tmp/cron-drupal-sync-w3-logs-"`bc <<< $(date +%j)%2`'.log'
su w3-drupal -c 'bin/sync-w3-logs' >> $log 2>&1
