# etc/cron.daily/update-node-hits
#!/bin/bash

cd /var/www/w3.uib.no

bin/site-drush scr update-node-hits
